# deploy/keylogme-zero.service
[Unit]
Description= Keylogme zero (keylogme.com)

[Install]
WantedBy=multi-user.target


[Service]
ExecStart=/bin/keylogme-zero
